# services/chatai-frontend.service.template
[Unit]\nDescription=NASP Chatbot Frontend\nAfter=network.target\n\n[Service]\nType=simple\nUser=dj\nWorkingDirectory=/home/dj/chatai\nEnvironment=PATH=/usr/local/bin:/usr/bin\nExecStart=/usr/bin/npm start\nRestart=always\nRestartSec=5\nStartLimitIntervalSec=0\n\nStandardOutput=append:/var/log/chatai-frontend.log\nStandardError=append:/var/log/chatai-frontend.error.log\n\n[Install]\nWantedBy=multi-user.target
